# systemd/rpi-backup-weekly.timer
[Unit]
Description=Weekly full Raspberry Pi image backup (when ISO week changes)

[Timer]
OnCalendar=Mon 00:00
Persistent=true
Unit=rpi-backup.service

[Install]
WantedBy=timers.target